# deb/install/prerm
#!/bin/bash

systemctl stop waggle-sanity-check.service || true
systemctl disable waggle-sanity-check.service || true
